# demo-secrets.conf
[program:demo-secrets]
environment=NODE_ENV="DEVELOPMENT"
command=/home/pi/.nvm/versions/node/v11.14.0/bin/node /home/pi/src/demo-secrets/app.js
directory=/home/pi/src/demo-secrets
